# deploy/corebit.service
[Unit]
Description=CoreBit Network Manager
Documentation=https://github.com/your-org/corebit
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=corebit
Group=corebit
WorkingDirectory=/opt/corebit
EnvironmentFile=/opt/corebit/.env
ExecStart=/usr/bin/node /opt/corebit/index.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=corebit

# Graceful shutdown timeout
TimeoutStopSec=30

# Network capabilities for device probing (ICMP ping, SNMP)
# Required on Ubuntu 24.04+ where setuid ping was removed
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN

# Security hardening (NoNewPrivileges=false required for AmbientCapabilities)
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/corebit/data /opt/corebit/backups /opt/corebit
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65535
MemoryMax=1G
CPUQuota=80%

# Network access required for device probing
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
